# config/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
